# scripts/launchd/com.fitme.hadf.phase2.collect.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  HADF Phase 2 cloud-fingerprinting collection schedule.

  Fires the collection wrapper 5 times per day. Each invocation makes
  50 calls per endpoint (openai, anthropic, local) and appends to
  .claude/shared/hadf/phase2-fingerprint-raw.jsonl.

  Hours below are LOCAL time. The wrapper tags each call with the
  UTC-derived window so post-hoc analysis can group correctly even if
  the local schedule does not fall exactly on the pre-registered UTC
  windows (02/08/14/18/22 UTC). The defaults here are spaced ~5h apart
  so a 3-day collection produces 5 x 3 = 15 runs x 150 calls = 2,250
  data points (well above the 600-point validity floor).

  Deploy:
    cp scripts/launchd/com.fitme.hadf.phase2.collect.plist \
       ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) \
       ~/Library/LaunchAgents/com.fitme.hadf.phase2.collect.plist
    launchctl print gui/$(id -u)/com.fitme.hadf.phase2.collect | head -40

  Stop / unload:
    launchctl bootout gui/$(id -u)/com.fitme.hadf.phase2.collect

  Adjust hours for your timezone if you want the windows closer to the
  pre-registered UTC targets. The wrapper's tagging logic does not care.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fitme.hadf.phase2.collect</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Volumes/DevSSD/FitTracker2/scripts/hadf-phase2-collect.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Volumes/DevSSD/FitTracker2</string>

    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key><integer>1</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>5</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>11</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>17</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>21</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Volumes/DevSSD/FitTracker2/.claude/shared/hadf/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Volumes/DevSSD/FitTracker2/.claude/shared/hadf/logs/launchd-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <!-- HADF_ENDPOINTS: skip 'local' for the 2026-04-30 campaign; Ollama on this MacBook Air
             runs llama3.2:3b at ~0.7 tps which is far below the harness's 60s urllib timeout.
             Pre-registration explicitly allows the 2-endpoint fallback path. Re-add 'local' if
             Ollama performance is fixed before campaign end. -->
        <key>HADF_ENDPOINTS</key>
        <string>openai,anthropic</string>
    </dict>
</dict>
</plist>
